# sample_config/pm.service
[Unit]
Description=gunicorn socket

[Socket]
ListenStream=/run/pm-gunicorn/socket
ListenStream=127.0.0.1:8000

[Install]
WantedBy=sockets.target


